# h1ds/deployment/management/upstart/h1ds-tasks.conf.template
description "h1ds-tasks"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid $user
setgid $group

exec ${deployment_dir}/tasks/tasks.sh
